# board/opendingux/target_skeleton/etc/profile.d/banner.sh
if [ "$PS1" ]; then
	cat /etc/welcome-banner.txt
fi
