# support/testing/tests/package/test_virglrenderer/rootfs-overlay/etc/profile.d/stty-raw.sh
# Avoid double-cooking the inner serial console, otherwise the test
# infrastructure cannot reliably retrieve return codes.
grep -Fq br-inner-virgl /proc/cmdline && stty raw
