# board/opendingux/target_skeleton/etc/inittab
# /etc/inittab
#
# Copyright (C) 2001 Erik Andersen <andersen@codepoet.org>
#
# Note: BusyBox init doesn't support runlevels.  The runlevels field is
# completely ignored by BusyBox init. If you want runlevels, use
# sysvinit.
#
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console
# runlevels == ignored
# action    == one of sysinit, respawn, askfirst, wait, and once
# process   == program to run

# Startup the system
::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mkdir -p /dev/pts /dev/shm
::sysinit:/bin/mount -a
::sysinit:/bin/mkdir -p /var/run
null::sysinit:/bin/ln -sf /proc/self/fd /dev/fd
null::sysinit:/bin/ln -sf /proc/self/fd/0 /dev/stdin
null::sysinit:/bin/ln -sf /proc/self/fd/1 /dev/stdout
null::sysinit:/bin/ln -sf /proc/self/fd/2 /dev/stderr

# Remount the main partition read-only, nodev, nosuid
::sysinit:/bin/mount -o remount,ro,nodev,nosuid /boot

# now run any rc scripts
::sysinit:/etc/init.d/rcS

::sysinit:/bin/echo -n 1 > /proc/sys/vm/overcommit_memory
::sysinit:/bin/echo -n 0 > /proc/sys/vm/mmap_min_addr
::sysinit:/bin/echo -n 0 > /proc/sys/vm/page-cluster
::sysinit:/bin/echo '3 4 1 3' > /proc/sys/kernel/printk

tty1::respawn:/usr/bin/nice --5 /usr/bin/sudo -u od /usr/sbin/frontend_start

# Stuff to do before rebooting
::shutdown:/etc/init.d/rcK
::shutdown:/sbin/swapoff -a
::shutdown:/bin/umount -a -r
