# board/opendingux/gcw0/overlay/etc/udev/rules.d/81-modules-unload.rules
# Unload the WiFi module
SUBSYSTEM=="net", ACTION=="remove", KERNEL=="wlan0", RUN+="/sbin/rmmod 8192cu 8188eu"
