# board/opendingux/gcw0/overlay/etc/network/if-up.d/avahi
#!/bin/sh

[ "$IFACE" != "lo" ] || exit 0

/usr/sbin/avahi-daemon -c || /usr/sbin/avahi-daemon -D
