# board/opendingux/target_skeleton/etc/network/if-up.d/xinetd
#!/bin/sh

[ "$IFACE" != "lo" ] || exit 0

/sbin/start-stop-daemon -o -S -q -x /usr/sbin/xinetd -- -pidfile /run/xinetd.pid
