# package/busybox/telnetd.service
[Unit]
Description=Telnetd Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/telnetd -F

[Install]
WantedBy=multi-user.target
