# board/yio-remote/rpi0/rootfs_overlay/etc/systemd/system/sharp-init.service
[Unit]
Description=Sharp init
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/yio-remote/sharp-init
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
