# board/yio-remote/rpi0/rootfs_overlay/etc/systemd/system/sharp-init.service
[Unit]
Description=Sharp init
DefaultDependencies=false

[Service]
Type=simple
ExecStart=/usr/bin/yio-remote/sharp-init
RemainAfterExit=yes
Restart=no

[Install]
WantedBy=sysinit.target
