# board/yio-remote/rpi0/rootfs_overlay/etc/systemd/system/update.service
[Unit]
Description=Updater

[Service]
Type=simple
ExecStart=/usr/bin/updater.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
