# support/testing/tests/package/test_xen/overlay/etc/profile.d/stty-raw.sh
# Avoid double-cooking the terminal, otherwise the test infrastructure
# would not be able to retrieve return codes properly.
stty raw
